# root/etc/s6-overlay/s6-rc.d/init-logseq-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# permissions
lsiown -R kasm-user:kasm-user \
    /config
lsiown -R kasm-user:kasm-user \
    /opt/logseq
